# plasma/etc/local.d/guest-vm-tools.start
#!/bin/bash

# Check if is a Virtualbox Machine.
vbox=$(lspci | grep -i virtualbox)

# Check if is a Vmware Machine.
vmware=$(lspci | grep -i vmware)

# Runlevel Services directory.
runlevels="/etc/runlevels/default"

# message if work.
sucessfull="\n >>> Service is already enabled! \n"

# if nothign to do.
nothing="\n >>> Nothing to do... \n"

# if not have value;
if [ -n "$vbox" ];
then
	file=$(ls $runlevels/virtualbox-guest-additions)

	if [ -n "$file" ];
	then
	    echo -e $sucessfull
	else
		rc-update add virtualbox-guest-additions default
	fi
elif [ -n "$vmware" ];
then
	file=$(ls $runlevels/vmware-tools)

	if [ -n "$file" ];
	then
	    echo -e $sucessfull
	else
		rc-update add vmware-tools default
	fi
else
	echo -e $nothing
fi
