# xfce/etc/local.d/system-install-fix.start
#!/bin/bash

install=/installation

portage=$(ls /var/git/meta-repo)


function packages_fix {

	# exec file script "packages-fix.sh".
	sh $install/scripts/packages-fix.sh

	# remove installation script.
	rm -rf $install

	# remove the script.
	rm /etc/local.d/system-install-fix.start

}

packages_fix
